# src/daemon/systemd/user/pulseaudio.service.in
[Unit]
Description=Sound Service

[Service]
ExecStart=@PA_BINARY@
Restart=always

[Install]
Also=pulseaudio.socket
WantedBy=default.target
