# src/daemon/systemd/system/pulseaudio.service.in
[Unit]
Description = PulseAudio system instance

[Service]
ExecStart = @PA_BINARY@ --system --daemonize=no --disallow-exit
Restart = on-failure

[Install]
Also = pulseaudio.socket
